# deploy/sujood-api.service
# systemd unit for the Sujood backend API.
# Install to /etc/systemd/system/sujood-api.service on the Ubuntu server.
#
#   sudo cp deploy/sujood-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now sujood-api
#   sudo systemctl status sujood-api
#   journalctl -u sujood-api -f      # live logs
#
# Assumes the repo is at /opt/sujood-mats and was built with `npm run build`.
# DATABASE_URL and CORS_ORIGIN are read from /opt/sujood-mats/.env by dotenv.

[Unit]
Description=Sujood API server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sujood-mats
ExecStart=/usr/bin/node dist/server.cjs
Environment=NODE_ENV=production
Restart=always
RestartSec=3
# Run as a non-root user that owns /opt/sujood-mats (create with: sudo adduser --system --group sujood)
User=sujood
Group=sujood

[Install]
WantedBy=multi-user.target
